# bot.conf
[program:telegrambot]
directory = /home/ubuntu/tgoyatbot
command = python3 app.py

autostart = true
autorestart = true
environment = HOME = "/home/ubuntu", USER = "ubuntu"

stderr_logfile =/home/ubuntu/tgoyatbot/logfile_err.log
stdout_logfile =/home/ubuntu/tgoyatbot/logfile.log